# server/astro-sync.service
[Unit]
Description=Astro FITS sync — astropc USB disk → NAS
Documentation=file:///home/pyl/Documents/astrobatch/server/astro-sync.sh
After=network-online.target mnt-astropc.automount mnt-nas-input.mount
Wants=network-online.target mnt-astropc.automount mnt-nas-input.mount

[Service]
Type=simple
User=pyl
Group=pyl
ExecStart=/bin/bash /home/pyl/Documents/astrobatch/server/astro-sync.sh
Restart=on-failure
RestartSec=30

StandardOutput=append:/home/pyl/Documents/astrobatch/server/logs/astro-sync.log
StandardError=append:/home/pyl/Documents/astrobatch/server/logs/astro-sync.log

[Install]
WantedBy=multi-user.target
